# .config/anacrontab
#WARNING: autogenerated with bin/generate_user_dirs.sh,
#change the .template file and rerun the script
1	10	clean_tmp.daily	/home/vsimone/bin/clean_my_temp_dirs.sh
@monthly 0	backup_cron	crontab -l > /home/vsimone/backups/crontab.dump
