# .config/anacrontab
@daily	10	clean_tmp.daily	/home/simone/bin/clean_my_temp_dirs.sh
